# scripts/com.sourcea.chat-unify-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.sourcea.chat-unify-server</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>REPLACE_WITH_SOURCEA/scripts/serve-chat-unify.sh</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <false/>
  <key>StandardOutPath</key>
  <string>REPLACE_WITH_HOME/.sina/chat-unify-server.log</string>
  <key>StandardErrorPath</key>
  <string>REPLACE_WITH_HOME/.sina/chat-unify-server.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>SINA_SOURCE_A</key>
    <string>REPLACE_WITH_SOURCEA</string>
    <key>CHAT_UNIFY_UI_VERSION</key>
    <string>3.4.0</string>
  </dict>
</dict>
</plist>
