# scripts/com.sourcea.s10-eternal-audit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sourcea.s10-eternal-audit</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>source ~/.zshrc 2>/dev/null; cd /Users/sinakazemnezhad/Desktop/SourceA; python3 scripts/s10_eternal_audit_loop_v1.py --daily --json >> /Users/sinakazemnezhad/.sina/s10-eternal-audit.log 2>&1</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/sinakazemnezhad/.npm-global/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/sinakazemnezhad</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/sinakazemnezhad/.sina/s10-eternal-audit.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/sinakazemnezhad/.sina/s10-eternal-audit-err.log</string>
</dict>
</plist>
